# wifi_plugins_init
#!/bin/sh /etc/rc.common

START=100

start() {
	cd / && ./wifi_plugins &
}